# systemd/user/dovetail.service.d/override.conf
[Service]
Environment="DOVETAIL_LUA_SHARE=/home/james/code/dovetail/builddir/share"
